# gitea/docker/docker/etc/s6/gitea/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/gitea > /dev/null
    exec chpst -u git env HOME=/data/git /app/gitea/gitea web
popd
