# gitea/docker/docker/etc/s6/syslogd/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /root > /dev/null
    exec busybox syslogd -nS -O-
popd
